# ansible/roles/common/templates/prometheus-app.service.j2
[Unit]
Description={{ service_description }}
After=syslog.target network.target

[Service]
Type=simple

ExecStart={{ java_path }} -jar {{ app_path }}

[Install]
WantedBy=multi-user.target